# dotfiles/user/.config/systemd/user/redshift.service
[Unit]
Description=Redshift - color temp adjustment.

[Service]
Type=simple
ExecStart=-%h/.bin/redshift

[Install]
WantedBy=default.target
